# scripts/com.steveonan.tradingview-cdp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.steveonan.tradingview-cdp</string>
  <key>LimitLoadToSessionType</key>
  <string>Aqua</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/steveonan/Documents/Crypto Trading/TradingViewMCP/tradingview-mcp/scripts/launch_tv.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/tradingview_cdp.launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/tradingview_cdp.launchd.err</string>
</dict>
</plist>
